# cron/clamav_scan.cron
0 13 * * 6 root /opt/clamav/scripts/clamav_scan.sh
