# cron/clamav_scan.cron
30 11 * * * root /opt/clamav/scripts/clamav_scan.sh

