# cron/clamav_update.cron
30 11 * * * root /opt/clamav/scripts/clamav_update.sh

