# cron/clamav_update.cron
# Mise à jour signatures toutes les 4 heures
0 8 * * 6 root /opt/clamav/scripts/clamav_update.sh

