# systemd/clamav-update.service
[Unit]
Description=ClamAV Update Service

[Service]
Type=oneshot
ExecStart=/opt/clamav/scripts/clamav_update.sh
