# systemd/clamav-update.timer
[Unit]
Description=ClamAV Update tous les dimanches à 01:00

[Timer]
OnCalendar=Sun 01:00
Persistent=true

[Install]
WantedBy=timers.target
