# systemd/clamav-scan.service
[Unit]
Description=ClamAV Scan Service

[Service]
Type=oneshot
ExecStart=/opt/clamav/scripts/clamav_scan.sh
