# systemd/clamav-scan.timer
[Unit]
Description=ClamAV Scan tous les dimanche à 01:00

[Timer]
OnCalendar=Sun 01:00
Persistent=true

[Install]
WantedBy=timers.target
